# repo/docker/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/harborbite/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/harborbite/php-fpm.log
stderr_logfile=/var/log/harborbite/php-fpm-error.log

[program:queue-worker]
command=php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/harborbite/queue-worker.log
stderr_logfile=/var/log/harborbite/queue-worker-error.log

[program:scheduler]
command=bash -c "while true; do php /var/www/html/artisan schedule:run --verbose --no-interaction >> /var/log/harborbite/scheduler.log 2>&1; sleep 60; done"
autostart=true
autorestart=true
stdout_logfile=/var/log/harborbite/scheduler-stdout.log
stderr_logfile=/var/log/harborbite/scheduler-error.log
